# src/systemd/cjet.service.in
[Unit]
Description="cjet service"
After = network.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cjet
Restart=always
Type=forking

[Install]
WantedBy=multi-user.target
