# src/systemd/cjet.service.in
[Unit]
Description="cjet service"
After = network.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cjet
Restart=always
Type=forking

ExecStartPost=/bin/cp @CMAKE_INSTALL_PREFIX@/bin/avahi/cjet.service /etc/avahi/services
ExecStopPost=/bin/rm /etc/avahi/services/cjet.service

[Install]
WantedBy=multi-user.target
