# src/systemd/cjet.service.in
[Unit]
Description="cjet service"
After = network.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cjet
Restart=always
Type=forking

ExecStartPost=/bin/cp @CMAKE_INSTALL_PREFIX@/bin/avahi/jet.service /etc/avahi/services
ExecStopPost=/bin/rm /etc/avahi/services/jet.service
ExecStartPost=/bin/cp @CMAKE_INSTALL_PREFIX@/bin/avahi/jetws.service /etc/avahi/services
ExecStopPost=/bin/rm /etc/avahi/services/jetws.service

[Install]
WantedBy=multi-user.target
